# deploy/cedgym-bot-watchdog.service
[Unit]
Description=CED-GYM WhatsApp Bot watchdog (auto-recovery)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cedgym-bot-watchdog.sh
# Hard timeout — si el script se cuelga (raro pero seguro defensivo)
TimeoutStartSec=300
# Nice level alto para no competir con Docker daemon ni el bot mismo
Nice=10
